# target/linux/mediatek/base-files/etc/uci-defaults/99_fwenv-store-ethaddr.sh
[ ! -e /etc/fw_env.config ] && exit 0

. /lib/functions/system.sh

exit 0
