# target/linux/mediatek/filogic/base-files/etc/hotplug.d/firmware/11-mt76-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"mediatek/mt7981_eeprom_mt7976_dbdc.bin")
	case "$board" in
	cmcc,rax3000m)
		case "$(cmdline_get_var root)" in
		/dev/mmc*)
			caldata_extract_mmc "factory" 0x0 0x1000
			;;
		esac
		;;
	esac
	;;
*)
	exit 1
	;;
esac
